# etc/init/samsara-cell-collector.conf
description "Samsara Cell Collector"
author "vilnei.neves@inf.ufpel.edu.br"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec samsara-cell-collector --config-file /etc/samsara/samsara.conf --log-file /var/log/samsara/samsara-cell-collector.log
